# supervisord.conf
[program:signald]
command=/usr/bin/signald -d /var/lib/signald --system-socket
environment=PATH="/usr/local/bin:/usr/bin:/bin"
user=signald
autostart=true
autorestart=true
redirect_stderr=true
